# pki/base/ocsp/shared/lib/systemd/system/pki-ocspd.target
[Unit]
Description=PKI Online Certificate Status Protocol Server
After=syslog.target network.target

[Install]
WantedBy=multi-user.target


