# base/ra/lib/systemd/system/pki-rad.target
[Unit]
Description=PKI Registration Authority Server
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
