# pki/base/tks/shared/lib/systemd/system/pki-tksd.target
[Unit]
Description=PKI Token Key Service
After=syslog.target network.target

[Install]
WantedBy=multi-user.target


