# base/server/share/lib/systemd/system/pki-tomcatd.target
[Unit]
Description=PKI Tomcat Server
Wants=dirsrv.target
After=syslog.target network.target dirsrv.target

[Install]
WantedBy=multi-user.target
